# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:deepflow-singularity-server]
command=node /app/dist/server.js
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=999
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=NODE_ENV="production"

[program:singularity-video-callback]
command=node /root/.openclaw/workspace-singularity-video/bin/video-agent-callback-server.mjs
directory=/root/.openclaw/workspace-singularity-video
autostart=true
autorestart=true
startsecs=3
startretries=999
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=NODE_ENV="production"
